# scripts/systemd/airflow-remove_worker_logs.timer
[Timer]
OnCalendar=*-*-* 00:00:00

[Install]
RequiredBy=airflow-workerset.service