# scripts/systemd/airflow-confapply-agent.path
[Unit]
After=airflow-scheduler.service airflow-webserver.service airflow-workerset.service airflow-workerset-small.service
PartOf=airflow-scheduler.service airflow-webserver.service airflow-workerset.service airflow-workerset-small.service

[Path]
PathModified=/etc/sysconfig/airflow.env

[Install]
WantedBy=airflow-scheduler.service airflow-webserver.service airflow-workerset.service airflow-workerset-small.service
