# scripts/systemd/airflow-heartbeat.timer
[Unit]
BindsTo=airflow-workerset.service

[Timer]
OnCalendar=*:0/1

[Install]
RequiredBy=airflow-workerset.service