# scripts/systemd/airflow-heartbeat.service
[Service]
EnvironmentFile=/etc/sysconfig/airflow.env
Type=oneshot
ExecStart=/opt/turbine/airflow-heartbeat.sh
